# etc/scripts/cloudhub.service
# If you modify this, please also make sure to edit init.sh

[Unit]
After=network-online.target

[Service]
User=snet
Group=snet
Environment="HOST=0.0.0.0"
Environment="PORT=8888"
#Environment="PORT=443"
#Environment="TLS_CERTIFICATE=/usr/lib/cloudhub/key/cloudhub_self_signed.pem"
Environment="BOLT_PATH=/var/lib/cloudhub/cloudhub-v1.db"
Environment="CANNED_PATH=/usr/share/cloudhub/cloudhub-canned"
Environment="PROTOBOARDS_PATH=/usr/share/cloudhub/cloudhub-protoboards"
Environment="TEMPLATES_PATH=/usr/share/cloudhub/cloudhub-templates"
EnvironmentFile=-/etc/default/cloudhub
ExecStart=/usr/bin/cloudhub $CLOUDHUB_OPTS
KillMode=control-group
Restart=on-failure

[Install]
WantedBy=multi-user.target